# supervisord.ini
[supervisord]
loglevel=warn
nodaemon=true
user=root

[program:jumpserver]
command=/bin/sh -c "/opt/jumpserver/jumpserver.sh"

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
