# supervisord.ini
[supervisord]
loglevel=warn
nodaemon=true
user=root
pidfile=/run/supervisord.pid

[program:jumpserver]
command=/bin/sh -c "/opt/jumpserver/jumpserver.sh"
redirect_stderr=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
redirect_stderr=true

