# supervisord.ini
[supervisord]
loglevel=warn
nodaemon=true
user=root
pidfile=/run/supervisord.pid

[program:jumpserver]
command=/usr/bin/python jms start all
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
