# supervisord/supervisord.ini
[supervisord]
loglevel=warn
nodaemon=true
user=root
pidfile=/run/supervisord.pid

[program:jumpserver]
command=/usr/local/bin/python jms start all
directory=/opt/jumpserver
exitcodes=0
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
